# packaging/systemd/okf-http.user.service.example
[Unit]
Description=OKF HTTP server
Documentation=https://github.com/kathrinmotzkus/open-knowledge-format
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/okf-http --browser-root /usr/share/okf/browser --host 127.0.0.1 8003
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=default.target
